# dot_config/systemd/user/docker_prune.timer
[Unit]
Description=Run docker prun every 24 hours

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
