# dot_config/systemd/user/vdirsyncer.service
[Unit]
Description=Run vdirsyncer

[Service]
ExecStart=/usr/bin/vdirsyncer

[Install]
WantedBy=default.target
